# example-setup/upstart.conf
# camus - deploy manager
#

description "camus"

start on filesystem
stop on runlevel [06]

respawn
setuid app
setgid app

script
    cd /home/app
    exec camus -server -enforce
end script
